# rc.local
#!/bin/bash

su -c /home/halloffame/hall-of-fame/launch.sh halloffame
